# docker/app/conf/cron.d/crontab
* * * * * www-data sh -c "php /app/bin/console cron:run" >> /app/var/log/cron.log
